# docker/development/etc/s6-overlay/s6-rc.d/laravel-scheduler/run
#!/command/with-contenv sh
# Drop privileges to www-data before exec — s6 starts as root
exec php /var/www/html/artisan schedule:work
