# docker/etc/s6-overlay/s6-rc.d/laravel-queue-2/run
#!/command/with-contenv sh
exec php /var/www/html/artisan queue:work \
    --tries=3 \
    --timeout=90 \
    --sleep=3 \
    --max-time=3600
