# docker/etc/s6-overlay/s6-rc.d/laravel-queue-1/finish
#!/command/with-contenv sh
# $1 = exit code of run script, $2 = signal that killed it
# Give queue:work a moment to finish the current job on SIGTERM
sleep 1
