# docker/production/etc/s6-overlay/s6-rc.d/zepeed-default-account/run
#!/command/with-contenv sh
# Runs once after S6 brings up the base services.
# Creates the default admin account if it does not already exist.
# Idempotent — safe on every container restart.
exec php /var/www/html/artisan app:create-user-account --default
