# docker/production/etc/s6-overlay/s6-rc.d/zepeed-default-account/up
/var/www/html/artisan app:create-user-account --default
