# docker/etc/s6-overlay/s6-rc.d/laravel-queue-4/run
#!/command/with-contenv sh
# Drop privileges to www-data before exec — s6 starts as root
exec php /var/www/html/artisan queue:work \
    --queue=realtime,mail,default \
    --tries=3 \
    --timeout=90 \
    --sleep=3 \
    --max-jobs=500
